# deploy/faceless-cartographer.service
[Unit]
Description=faceless-cartographer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Adjust User/WorkingDirectory to match the account/checkout on the Pi.
# PORT (and anything else) is read from .env in WorkingDirectory (`just serve`
# passes --env-file-if-exists=.env to tsx/node) -- not set here, so the same
# .env used for manual `just serve` runs is also what systemd uses.
User=pi
WorkingDirectory=/home/pi/faceless-cartographer
Environment=NODE_ENV=production
# Absolute path: systemd units don't get a login shell's PATH. Adjust if
# `just` was installed somewhere other than /usr/local/bin (check with
# `which just`).
ExecStart=/usr/local/bin/just serve
Restart=on-failure
RestartSec=5
# Hardening (safe defaults for a small single-purpose service). If the D8
# LLM explanation feature (CLI call to `claude -p`) is implemented later,
# revisit this: the `claude` CLI may need read access to its own config/auth
# directory (e.g. ~/.claude), which ProtectSystem=strict would block.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/pi/faceless-cartographer

[Install]
WantedBy=multi-user.target
